# packages/010_callback_phishing_rmm_persistence_defense/detections/crowdstrike/callback_phishing_rmm_ioa.json
{"ioa_definitions":[{"name":"Callback Phishing - RMM Installation from Browser","description":"Detects RMM tools downloaded and executed directly from browser, indicating callback phishing success","severity":"critical","platform":"windows","tactic":"Command and Control","technique":"T1219","disposition":"BLOCK_AND_ALERT","pattern":{"type":"process_hierarchy","conditions":{"parent_process":["chrome.exe","msedge.exe","firefox.exe"],"child_process":["AnyDesk.exe","TeamViewer.exe","ScreenConnect.exe"]}},"response_actions":["kill_process","isolate_host","alert_soc"]},{"name":"Unauthorized RMM Persistence","description":"Detects creation of services or scheduled tasks for unauthorized RMM tools","severity":"high","platform":"windows","tactic":"Persistence","technique":"T1543.003","disposition":"BLOCK_AND_ALERT","pattern":{"type":"service_creation","conditions":{"service_name_pattern":["*AnyDesk*","*TeamViewer*","*ScreenConnect*","*ConnectWise*"],"exclude_approved_paths":true}},"response_actions":["block_service","remove_persistence","alert_soc"]}],"metadata":{"package":"pkg-010","version":"1.0.0","last_updated":"2025-03-28"}}
