# packages/019_storm_1175_high_velocity_medusa_ransomware_campaign_defense/hunts/analysis/019_storm_1175_high_velocity_medusa_ransomware_campaign_defense_hunt_procedure.yml
apiVersion: soac.io/v1
kind: ThreatHunt
metadata:
  name: storm-1175-medusa-ransomware-defense
  version: 1.0.0
  author: SOaC Hunt Engineering Team
  created: '2026-04-13T20:14:47.849Z'
  package_id: pkg-019
  mitre_attack:
    - T1190
    - T1059
    - T1053
    - T1003
    - T1552
    - T1021
    - T1570
    - T1005
    - T1486
  severity: CRITICAL
  data_sources:
    - MDE
    - M365
    - CrowdStrike
    - Sentinel
    - Entra ID
    - AWS
  tags:
    - ransomware
    - medusa
    - storm-1175
    - n-day
    - exploitation
    - lateral-movement
spec:
  hypothesis: >-
    Storm-1175 threat actors are exploiting known vulnerabilities in
    internet-facing applications to deploy Medusa ransomware through rapid
    post-exploitation techniques
  hunt_type: proactive
  analysis_steps:
    - id: step-1
      name: Initial Exploitation Validation
      description: >-
        Analyze query-001 results to identify compromised internet-facing
        applications and validate initial access vectors
      tools:
        - Microsoft Defender for Endpoint
        - Process Explorer
        - Timeline Analysis
    - id: step-2
      name: Credential Compromise Assessment
      description: >-
        Review query-002 findings to determine scope of credential harvesting
        and identify compromised accounts
      tools:
        - Azure AD Sign-in Logs
        - Bloodhound
        - Account Activity Analysis
    - id: step-3
      name: Lateral Movement Mapping
      description: >-
        Correlate query-003 and query-004 results to map lateral movement paths
        and persistence mechanisms
      tools:
        - Network Topology Mapping
        - Authentication Log Analysis
        - Task Scheduler Investigation
    - id: step-4
      name: Impact Assessment and Containment
      description: >-
        Evaluate query-005 results to assess ransomware deployment scope and
        initiate containment procedures
      tools:
        - File System Analysis
        - Network Segmentation
        - Backup Verification
  success_criteria:
    - >-
      Identification of exploitation attempts against internet-facing
      applications with subsequent command execution
    - >-
      Detection of credential harvesting activities targeting LSASS or using
      known dumping techniques
    - >-
      Evidence of lateral movement using compromised credentials across multiple
      systems
    - >-
      Discovery of scheduled tasks or other persistence mechanisms created by
      suspicious processes
    - >-
      Confirmation of file encryption patterns consistent with Medusa ransomware
      deployment
