# packages/015_apt37_famous_chollima_operation_hankook_phantom_defense/scenario.json
{
  "package_id": "pkg-015",
  "title": "APT37 (Famous Chollima) - Operation HanKook Phantom Defense",
  "version": "1.0.0",
  "steps": [
    {
      "step_index": 1,
      "role": "Red Team",
      "phase": "Initial Access",
      "pillar": "Offense",
      "title": "Spear-phishing with malicious decoy documents",
      "text": "APT37 sends targeted spear-phishing emails containing malicious decoy documents (HWP, PDF, LNK files) such as 'National Intelligence Research Society Newsletter' to South Korean academics, think tanks, and energy/security associations",
      "description": "APT37 sends targeted spear-phishing emails containing malicious decoy documents (HWP, PDF, LNK files) such as 'National Intelligence Research Society Newsletter' to South Korean academics, think tanks, and energy/security associations",
      "mitre": [
        "T1566.001",
        "T1566.002"
      ],
      "expectedOutcome": "Step 1 completes successfully",
      "detectionHints": "Monitor for T1566.001, T1566.002 indicators",
      "playbookRefs": []
    },
    {
      "step_index": 2,
      "role": "Red Team",
      "phase": "Initial Access",
      "pillar": "Offense",
      "title": "Initial payload execution via LNK files",
      "text": "Malicious LNK files are executed on victim systems, initiating the infection chain and preparing for malware deployment",
      "description": "Malicious LNK files are executed on victim systems, initiating the infection chain and preparing for malware deployment",
      "mitre": [
        "T1204.002",
        "T1059.001"
      ],
      "expectedOutcome": "Step 2 completes successfully",
      "detectionHints": "Monitor for T1204.002, T1059.001 indicators",
      "playbookRefs": []
    },
    {
      "step_index": 3,
      "role": "Red Team",
      "phase": "Execution",
      "pillar": "Offense",
      "title": "PowerShell-based infection chain deployment",
      "text": "The attack leverages PowerShell scripts to execute the infection chain, downloading and installing the RokRat backdoor on compromised systems",
      "description": "The attack leverages PowerShell scripts to execute the infection chain, downloading and installing the RokRat backdoor on compromised systems",
      "mitre": [
        "T1059.001",
        "T1105"
      ],
      "expectedOutcome": "Step 3 completes successfully",
      "detectionHints": "Monitor for T1059.001, T1105 indicators",
      "playbookRefs": []
    },
    {
      "step_index": 4,
      "role": "Red Team",
      "phase": "Execution",
      "pillar": "Offense",
      "title": "RokRat backdoor installation and persistence",
      "text": "RokRat malware is installed with persistence mechanisms including scheduled tasks, registry modifications, and remote management tools to maintain long-term access",
      "description": "RokRat malware is installed with persistence mechanisms including scheduled tasks, registry modifications, and remote management tools to maintain long-term access",
      "mitre": [
        "T1053.005",
        "T1112",
        "T1219"
      ],
      "expectedOutcome": "Step 4 completes successfully",
      "detectionHints": "Monitor for T1053.005, T1112, T1219 indicators",
      "playbookRefs": []
    },
    {
      "step_index": 5,
      "role": "SOC Team",
      "phase": "Command and Control",
      "pillar": "Respond",
      "title": "Cloud-based C2 communication establishment",
      "text": "The malware establishes encrypted command and control communications using compromised cloud services including Dropbox, Yandex, and pCloud for data exfiltration and remote control",
      "description": "The malware establishes encrypted command and control communications using compromised cloud services including Dropbox, Yandex, and pCloud for data exfiltration and remote control",
      "mitre": [
        "T1071.001",
        "T1567.002",
        "T1573.001"
      ],
      "expectedOutcome": "Step 5 completes successfully",
      "detectionHints": "Monitor for T1071.001, T1567.002, T1573.001 indicators",
      "playbookRefs": []
    },
    {
      "step_index": 6,
      "role": "SOC Team",
      "phase": "Exfiltration",
      "pillar": "Respond",
      "title": "Data collection and exfiltration",
      "text": "Sensitive data is collected from compromised systems and exfiltrated through the established cloud-based C2 infrastructure to support DPRK espionage objectives",
      "description": "Sensitive data is collected from compromised systems and exfiltrated through the established cloud-based C2 infrastructure to support DPRK espionage objectives",
      "mitre": [
        "T1005",
        "T1041",
        "T1567.002"
      ],
      "expectedOutcome": "Step 6 completes successfully",
      "detectionHints": "Monitor for T1005, T1041, T1567.002 indicators",
      "playbookRefs": []
    }
  ]
}