# packages/010_callback_phishing_rmm_persistence_defense/mappings/mitre_attack.json
{"package_id":"pkg-010","framework":"MITRE ATT&CK v14","techniques":[{"id":"T1566.003","name":"Phishing: Spearphishing via Service","tactic":"Initial Access","detections":["callback_phishing_detection","callback_phishing_chain"]},{"id":"T1219","name":"Remote Access Software","tactic":"Command and Control","detections":["rmm_tool_abuse","rmm_abuse_detection","rmm_c2_network_detection"]},{"id":"T1543.003","name":"Windows Service","tactic":"Persistence","detections":["rmm_persistence_mechanisms"]},{"id":"T1053.005","name":"Scheduled Task","tactic":"Persistence","detections":["rmm_persistence_mechanisms"]},{"id":"T1204.001","name":"User Execution: Malicious Link","tactic":"Execution","detections":["rmm_tool_abuse"]}]}
