# packages/014_axios_supply_chain_compromise_defense/detection.yaml
apiVersion: soac.io/v1
kind: DetectionRule
metadata:
  name: Supply Chain Package Compromise Defense
  version: "1.0"
  author: SOaC Detection Engineering Team
  created: "2024-01-15"
  package_id: "pkg-014"
  mitre_attack:
    - T1195.002
    - T1204.005
    - T1027
    - T1082
    - T1071.001
    - T1105
    - T1036.005
    - T1547.001
    - T1059.001
    - T1059.006
    - T1219
    - T1059
    - T1083
    - T1070.004
    - T1070.006
  severity: "critical"
  data_sources:
    - process_creation
    - network_traffic
    - file_creation
    - registry_modification
    - command_execution
spec:
  description: "Detects supply chain package compromise attacks targeting npm packages with malicious lifecycle hooks, C2 communication, RAT deployment, and anti-forensic activities commonly associated with DPRK threat actors"
  platforms:
    - windows
    - linux
    - macos
  inline_rules:
    - id: "rule-014-001"
      name: "Suspicious NPM Package Installation with Postinstall Scripts"
      format: "kql"
      platform: "windows"
      mitre:
        - T1195.002
        - T1204.005
      severity: "high"
      query: |
        DeviceProcessEvents
        | where ProcessCommandLine contains "npm" and ProcessCommandLine contains "install"
        | where ProcessCommandLine contains_any ("postinstall", "preinstall", "install")
        | where ProcessCommandLine contains_any ("--unsafe-perm", "--allow-root", "--ignore-scripts=false")
        | where ParentProcessName in~ ("cmd.exe", "powershell.exe", "node.exe")
        | project Timestamp, DeviceName, ProcessCommandLine, ParentProcessName, InitiatingProcessAccountName
    
    - id: "rule-014-002"
      name: "Obfuscated Script Execution from Package Directories"
      format: "sigma"
      platform: "linux"
      mitre:
        - T1027
        - T1059.001
      severity: "high"
      query: |
        title: Obfuscated Script Execution from Package Directories
        logsource:
          category: process_creation
          product: linux
        detection:
          selection:
            Image|endswith:
              - '/node'
              - '/bash'
              - '/sh'
            CommandLine|contains:
              - 'node_modules'
              - '.npm'
            CommandLine|contains:
              - 'eval('
              - 'Buffer.from('
              - 'atob('
              - 'base64'
              - 'hex'
          condition: selection
        falsepositives:
          - Legitimate package installations
        level: high
    
    - id: "rule-014-003"
      name: "Suspicious Network Communication from Node Processes"
      format: "splunk"
      platform: "windows"
      mitre:
        - T1071.001
        - T1082
      severity: "medium"
      query: |
        index=network sourcetype=firewall OR sourcetype=proxy
        | where process_name="node.exe" OR parent_process="npm.exe"
        | where dest_port NOT IN (80,443,22,53)
        | where dest_ip NOT IN (registry.npmjs.org, github.com, raw.githubusercontent.com)
        | stats count by src_ip, dest_ip, dest_port, process_name
        | where count > 5
    
    - id: "rule-014-004"
      name: "RAT Payload Deployment and Masquerading"
      format: "kql"
      platform: "windows"
      mitre:
        - T1105
        - T1036.005
      severity: "critical"
      query: |
        DeviceFileEvents
        | where FolderPath contains_any ("node_modules", "npm", "temp", "tmp")
        | where FileName matches regex @"^[a-zA-Z0-9]{8,16}\.(exe|dll|bat|ps1)$"
        | where ActionType == "FileCreated"
        | join kind=inner (
            DeviceNetworkEvents
            | where RemoteUrl !contains_any ("npmjs.org", "github.com", "microsoft.com")
        ) on DeviceId
        | project Timestamp, DeviceName, FolderPath, FileName, RemoteUrl, InitiatingProcessCommandLine
    
    - id: "rule-014-005"
      name: "Cross-Platform Persistence Establishment"
      format: "sigma"
      platform: "macos"
      mitre:
        - T1547.001
        - T1036.005
      severity: "high"
      query: |
        title: Cross-Platform Persistence Establishment
        logsource:
          category: file_event
          product: macos
        detection:
          selection:
            TargetFilename|contains:
              - '/Library/LaunchAgents/'
              - '/Library/LaunchDaemons/'
              - '$HOME/Library/LaunchAgents/'
            TargetFilename|endswith: '.plist'
            ProcessImage|contains:
              - 'node'
              - 'npm'
          condition: selection
        falsepositives:
          - Legitimate application installations
        level: high
    
    - id: "rule-014-006"
      name: "Remote Access Tool Command Execution"
      format: "wazuh"
      platform: "linux"
      mitre:
        - T1219
        - T1059
      severity: "critical"
      query: |
        <rule id="100014" level="12">
          <if_group>syscheck</if_group>
          <field name="file">node_modules|\.npm</field>
          <field name="command">sh|bash|curl|wget</field>
          <regex>eval|exec|system</regex>
          <description>Remote access tool command execution from package directories</description>
          <mitre>
            <id>T1219</id>
            <id>T1059</id>
          </mitre>
        </rule>
    
    - id: "rule-014-007"
      name: "Anti-Forensic File and Log Cleanup"
      format: "splunk"
      platform: "windows"
      mitre:
        - T1070.004
        - T1070.006
      severity: "medium"
      query: |
        index=windows EventCode=4663 OR EventCode=4656
        | where Object_Name LIKE "%npm-debug.log%" OR Object_Name LIKE "%node_modules%"
        | where Access_Mask="0x10000" OR Access_Mask="DELETE"
        | where Process_Name!="npm.exe" AND Process_Name!="node.exe"
        | stats count by Computer_Name, Process_Name, Object_Name, Account_Name
        | where count > 3
  
  test_cases:
    - name: "Malicious NPM Package Installation"
      description: "Simulates installation of compromised npm package with suspicious postinstall script"
      attack_techniques:
        - T1195.002
        - T1204.005
      test_data:
        process_command: "npm install malicious-package --unsafe-perm"
        parent_process: "cmd.exe"
        file_path: "C:\\Users\\victim\\AppData\\Roaming\\npm\\node_modules\\malicious-package"
      expected_detection: true
    
    - name: "Obfuscated Payload Execution"
      description: "Tests detection of base64 encoded script execution from node_modules directory"
      attack_techniques:
        - T1027
        - T1059.001
      test_data:
        process_command: "node -e \"eval(Buffer.from('Y29uc29sZS5sb2coJ21hbGljaW91cycpOw==', 'base64').toString())\""
        working_directory: "/home/user/project/node_modules/.bin"
        network_connection: "192.168.1.100:4444"
      expected_detection: true
    
    - name: "Cross-Platform RAT Persistence"
      description: "Validates detection of persistence mechanism creation by compromised package"
      attack_techniques:
        - T1547.001
        - T1036.005
        - T1105
      test_data:
        file_created: "/Library/LaunchAgents/com.apple.update.plist"
        process_name: "node"
        network_download: "http://malicious-c2.com/payload.bin"
        masquerading_name: "system_update"
      expected_detection: true